# ansible/files/api.service
[Unit]
Description=App

[Service]
ExecStart=/home/ubuntu/app

[Install]
WantedBy=multi-user.target