# ansible/files/app.service
[Unit]
Description=App

[Service]
ExecStart=/home/ubuntu/bin/app

[Install]
WantedBy=multi-user.target